# Client/supervisord_prod.conf
[program:pyplc]
directory=/home/pi/pythonPLC/Client/
command=/home/pi/.pyenv/versions/pyplc/bin/python /home/pi/pythonPLC/Client/manage.py --config prod --url server --start
stopsignal=KILL
stopasgroup=true
autostart=true
autorestart=true
startsecs=15
# startretries=36
stopwaitsecs=10
stdout_logfile=/home/pi/pyplc.log
stderr_logfile=/home/pi/pyplc.err
